# root/etc/services.d/ffmpeg/run
#!/usr/bin/with-contenv bash

echo "Starting ffmpegMovie."
cd /

exec \
s6-setuidgid abc bash ffmpegMovie